# systemd/.config/systemd/user/low-battery.timer
[Unit]
Description=Check battery level

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
